# roles/domoticz/templates/domoticz-get-update.timer.j2
[Unit]
Description=Update Domoticz daily
Before=domoticz.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
