# roles/nginx/templates/nginx.socket.j2
[Unit]
After=network.target

[Socket]
ListenStream=80
ListenStream=0.0.0.0:80
BindIPv6Only=ipv6-only

[Install]
WantedBy=sockets.target
